# service/conf/nginx.ini
[program:nginx]
command=/opt/openresty/nginx/sbin/nginx -p /var/www/ -c /opt/openresty/nginx/conf/nginx.conf -g 'daemon off;'
process_name=%(program_name)s
priority=1
user=root
autostart=true
autorestart=true
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=5
serverurl=AUTO
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
