# service/conf/php-fpm.ini
[program:php-fpm]
command=/opt/php7/sbin/php-fpm -F -c /opt/php7/etc/php.ini
process_name=%(program_name)s
umask=022
user=root
autostart=true
exitcodes=0,2
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
